# systemd/phonchain-node.service
[Unit]
Description=Phonchain Node (PoP-S4.1)
After=network.target

[Service]
Type=simple
User=phonchain
Group=phonchain
WorkingDirectory=/opt/phonchain-node
EnvironmentFile=/opt/phonchain-node/.env
ExecStart=/opt/phonchain-node/.venv/bin/gunicorn -w 2 -b 127.0.0.1:5000 app:app
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
